# com.bug-tracker-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bug-tracker-mcp</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/arobel/repos/mcp/bug-tracker-mcp/.venv/bin/uv</string>
    <string>run</string>
    <string>/Users/arobel/repos/mcp/bug-tracker-mcp/server.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/arobel/repos/mcp/bug-tracker-mcp</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>OBSIDIAN_VAULT_PATH</key>
    <string>/Users/arobel/Obsidian/ND</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>30</integer>
  <key>StandardOutPath</key>
  <string>/tmp/bug-tracker-mcp.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/bug-tracker-mcp.err.log</string>
</dict>
</plist>
